# arm-ripper-build.timer
[Unit]
Description=My Scheduled Service Timer

[Timer]
OnCalendar=Sat *-*-* 21:00:00
Persistent=true

[Install]
WantedBy=timers.target